# hhkb/50-hhkb-add.rules
ACTION=="add", SUBSYSTEM=="hid", ENV{HID_NAME}=="HHKB-Hybrid_1" RUN+="/usr/local/bin/hhkb.sh %p"
